# files/etc/uci-defaults/commotion-quick-start
#!/bin/sh
[ -f "/usr/bin/patch" ] && { \
  cd ...controller...
  patch -p0 < /usr/share/commotion/patches/splash.qs.patch
}
